# bash/profile.d/aliases
#!/usr/bin/env bash

# Enable capslock as backspace and double shift as capslock.
alias k="setxkbmap -option caps:backspace && setxkbmap -option shift:both_capslock"

# print alias neatly.
alias path='printf "%b\n" "${PATH//:/\\n}"'

# System commands
alias brestart="sudo /etc/init.d/bluetooth restart"
