# bash/profile.d/aliases
#!/usr/bin/env bash

alias k="setxkbmap -option caps:backspace && setxkbmap -option shift:both_capslock | xset r 66"
alias K="k"

#------------------------------------------
#             Custom commands             |
#------------------------------------------

alias path='printf "%b\n" "${PATH//:/\\n}"'          # Print PATH env neatly.
alias bashopts='printf "%b\n" "${BASHOPTS//:/\\n}"'  # Show bashopts env neatly.
alias brestart="sudo /etc/init.d/bluetooth restart"  # Restart bluetooth.
alias mom="ps --no-headers -o comm $PPID"      # which process spawns the shell?
alias todo="grep --recursive 'TOD[O]: '"         # List all todo comments.
alias cwd='printf '%s' "${PWD##*/}"'         # Print current working directory.

#------------------------------------------
#              Git commands               |
#------------------------------------------

alias gl="git log --pretty=oneline --graph --abbrev-commit"
alias gs="git status --short"
alias gshow="git show"
alias gdf="git diff"
alias gca="git commit --all"

# Reset uncommited changes.
alias grest="git reset --hard HEAD"

#------------------------------------------
#             System commands             |
#------------------------------------------

alias sudo='sudo '                     # Allow aliases to be `sudo`ed.
#alias vim='$EDITOR '                  # Expand from default editor.
alias rm='rm '                         # Allow aliases to be `rm`ed.
alias vifm='vifm '                     # Allow aliase expansion for vifm paths.

# Sort, list and show hiden files/dirs with colored listing.
alias ls='lsc --almost-all --classify --group-directories-first '
alias dir='dr --almost-all --classify --group-directories-first '

#------------------------------------------
#     Files and directories shortcuts     |
#------------------------------------------

# Frequently used directories.
alias home="~"
alias p="~/Projects"
alias dot="~/Projects/dotfiles"
alias dt="~/Desktop"
alias m="~/Music"
alias dl="~/Downloads"
alias d="~/Documents"
alias et="/etc"

# Neovim Files.
alias lvim="~/.config/nvim/init.vim"          # Local nvim config.
alias svim="/usr/share/nvim/sysinit.vim"      # System-wide nvim config.
