# bash/profile.d/exports
#!/usr/bin/env bash

# Set `vim` as the default editor.
export EDITOR="vim"

# Set default term.
export TERM="xterm-256color"
