# bash/profile.d/exports
#!/usr/bin/env bash

# Set Neovim as the default editor.
export EDITOR='nvim '

# Neovim ENV, used to locate various files, such as the user's init.vim.
#export VIM="/usr/share/$EDITOR"

# Set default term.
export TERM="xterm-256color"

 # TODO: This should have a environmental path.
# Define local inputrc file.
export INPUTRC="/home/devs/Projects/dotfiles/bash/inputrc"

# Dotfiles executables.
export PATH="$PATH:/home/devs/Projects/dotfiles/bin"
