# bash/profile.d/exports
#!/usr/bin/env bash

# Set Neovim as the default editor.
export EDITOR='vim '

# --------------------------------------
#          FZF Default Options.         |
# --------------------------------------
export FZF_DEFAULT_OPTS='--no-reverse'       # Prompt always at the bottom.
export FZF_DEFAULT_OPTS+=' --border'         # Show top and bottom border.
export FZF_DEFAULT_OPTS+=' --no-bold'        # Disable bold characters.

export FZF_TMUX=1                     # Open in split-pane if tmux is running.

# Set default term.
export TERM="xterm-256color"

 # TODO: This should have a environmental path.
# Define local inputrc file.
export INPUTRC="/home/devs/Projects/dotfiles/bash/inputrc"


# --------------------------------------
#           Environment PATH.           |
# --------------------------------------

# Dotfiles executables.
if ! printenv PATH | grep 'dotfiles' &> /dev/null; then
  export PATH="$HOME/Projects/dotfiles/bin:$PATH"
fi
