# bash/profile.d/exports
#!/usr/bin/env bash

  # Set the default editor.
  EDITOR="$(which nvim || which vim)"

  # Increase bash history line limit (default 500).
  HISTFILESIZE=10000

  # Increase bash command history limit (default 500).
  HISTSIZE=10000

  HISTCONTROL=ignoreboth:erasedups
  # Ignore duplicate commands and spaces.

 # One history for all.
  PROMPT_COMMAND="history -a; $PROMPT_COMMAND"

  # Don't add this comman to history file.
  HISTIGNORE="bash:exit:pwd:dir:dirs:ls:clear:echo:man *"

  # Edit and run last command in history.
  FCEDIT="$EDITOR -u /home/devs/Projects/dotfiles/nvim/fcedit.vim"

  TERM="xterm-256color"

  # Custom executables.
  if ! printenv PATH | grep 'dotfiles' &> /dev/null; then
    PATH="$HOME/Projects/dotfiles/bin:$PATH"
  fi

  # fzf
  export FZF_DEFAULT_OPTS="\
     --height='40%' \
     --ansi    \
     --multi   \
     --reverse \
     --inline-info \
     --header='' \
     --color=dark,bg+:23,prompt:30,pointer:108 \
     --margin='1,1,1%' \
     --prompt='  "$(printf %s "${PWD##*/}")"  '
  "
