# bash/profile.d/exports
#!/usr/bin/env bash

# Set `vim` as the default editor.
export EDITOR="vim"
