# bash/profile.d/functions
#!/usr/bin/env bash

# Create directory and enter it.
function md() {

  local dir="${1}"

  mkdir "${dir}" && cd "${dir}" || return

}

# Colorized `cat` output with linenumbers in `less` pager.
function hicat() {

  local file="${1}"
  local filesize=""

  filesize="$(pygmentize -g "${file}" | wc -l)"

  if [ "$filesize" -gt "$(tput lines | awk '{print $file-3}')" ]; then
    pygmentize -g "${file}" | cat --number | less --RAW-CONTROL-CHARS
  else
    pygmentize -g "${file}" | cat --number
  fi

}

# Cd into the current directory after quiting vifm.
vifm() {

  local dest=""
  dest="$(command vifm --choose-dir -)"

  if [ -z "$dest" ]; then
    echo 'Directory picking cancelled/failed'
    return 1
  fi

  cd "$dest" || return
}
