# bash/profile.d/functions
#!/usr/bin/env bash

# Create directory and enter it.
function md() {

  local dir="${1}"

  mkdir "${dir}" && cd "${dir}" || return

}
