# bash/profile.d/options
#!/usr/bin/env bash

# automaticall preprend `cd` to directory names.
shopt -s autocd
